# homework.cron
*/20 6-9 * * 1-5 python3 /home/Fred/han/homework.py
